# mapcrafter-cron/crontab
# min   hour    day     month   weekday command
*/15    *       *       *       *       /generate_map.sh
